# scripts/voltaire-dispatch.service
[Unit]
Description=Voltaire Dispatch Service
After=network.target openclaw.service

[Service]
Type=simple
User=voltaire
WorkingDirectory=/opt/voltaire/dispatch-service
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=5
EnvironmentFile=/opt/voltaire/.env
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
